# images/lxc/unique-id-generator/scripts/unique-id-generator.service
[Unit]
Description=Unique ID Generator Service
After=network.target

[Service]
Type=simple
User=nobody
Group=nogroup
WorkingDirectory=/opt/unique-id-generator
ExecStart=/usr/bin/node /opt/unique-id-generator/server.js
Restart=always
RestartSec=10
StandardOutput=append:/var/log/unique-id-generator.log
StandardError=append:/var/log/unique-id-generator.log
Environment="NODE_ENV=production"
Environment="PORT=7001"
# Sharding configuration - will be overridden by container config
Environment="SHARD_ID=1"
Environment="TOTAL_SHARDS=1"

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/unique-id-generator /var/log

# Resource limits
MemoryHigh=200M
MemoryMax=256M
CPUQuota=50%

[Install]
WantedBy=multi-user.target